# deploy/upstart/logfiles.conf
# logfiles node
#
# this isn't run as a regular daemon but instead is started in screen
# so all we want from upstart is to fire our script at boot and then forget
# about us. that means this task is really "start logfiles node" rather
# than "logfiles node"

description "Start Logfiles Node"
author "Gavin Norman <gavin.norman@sociomantic.com>"

# we should start when the networking scripts have stopped (as presumably they
# will have brought the network up) and screen-cleanup is done so it is safe
# for us to start a new screen
start on static-network-up and stopped screen-cleanup

# we are a one off, non long running process
task

# title of screen session to create
env SCREEN_TITLE="logfiles"
# command we wish to run
env COMMAND="dhtnode-64"
# directories we wish to run $COMMAND in
env DIRS="/srv/dht/log-dht-1 /srv/dht/log-dht-2"

# set resource limits
limit nofile 100000 100000
limit core unlimited unlimited

# this is the actual command we want to run
exec /usr/local/bin/start_in_screen.sh -c "${COMMAND}" -t "${SCREEN_TITLE}" ${DIRS}
